# config/network/systemd/tatbot-network-auto-renew.service
[Unit]
Description=Tatbot Network Auto-Renewal Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/home/ook/tatbot/scripts/network_auto_renew.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target