# config/network/systemd/tatbot-mode-auto.service
[Unit]
Description=Tatbot Auto Mode Detection (Edge/Home)
After=network-online.target dnsmasq.service
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=root
ExecStart=/home/rpi2/tatbot/scripts/mode_auto_switcher_with_dhcp.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target


